# files/backup/mysql_backup.cron
00 01 * * * root /usr/bin/mysqldump --default-character-set=utf8 --all-databases --all --flush-logs --lock-tables --single-transaction | gzip > /var/lib/mysql/data/mysqldump.sql.gz && chmod 600 /var/lib/mysql/data/mysqldump.sql.gz
